# .config/autostart/auto-startup-workspace.desktop
[Desktop Entry]
Type=Application
Name=Auto Startup Workspace Manager
Comment=Automatically launch applications and assign them to workspaces
Exec=/home/bart/bin/auto-startup-workspace
Icon=preferences-system-startup
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
StartupNotify=false
